# deploy/etc/cron.d/restart_reshuffle_openttd
30 3 * * 0 root /usr/bin/systemctl restart openttd-dedicated.service
