# deploy/etc/cron.d/update_motd
* * * * * root /usr/local/bin/update_motd.sh
